# etc/systemd/system/zfs-scrub-weekly.service.d/timeout-example.conf
[Service]
# Override example: extend scrub timeout for large pools (>100TB)
# Copy this file as timeout.conf and adjust value per environment.
# 24h timeout:
TimeoutStartSec=86400
