# etc/systemd/firstboot/debvisor-firstboot.service
#
# DebVisor First Boot Initialization Service Unit
#
# Repository path: /etc/systemd/firstboot/debvisor-firstboot.service
# Installation target: /etc/systemd/system/debvisor-firstboot.service
#
# This unit must be installed to /etc/systemd/system/ for systemd auto-discovery.
# The installer or packaging script should copy/symlink this file there.
#
# Enable with: systemctl enable debvisor-firstboot.service
# Manual run: systemctl start debvisor-firstboot.service
#

[Unit]
Description=DebVisor First Boot Initialization
After=network-online.target
Wants=network-online.target
ConditionFirstBoot=yes
ConditionPathExists=!/var/lib/debvisor/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/debvisor-firstboot.sh
RemainAfterExit=yes
TimeoutSec=3600
Restart=on-failure
RestartSec=10s
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
