# opt/config/includes.chroot/etc/systemd/system/hostname-register.service
[Unit]
Description=DebVisor Hostname Registration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hostname-register.sh
RemainAfterExit=yes
SyslogIdentifier=hostname-register
ConditionPathIsExecutable=/usr/local/bin/hostname-register.sh

[Install]
WantedBy=multi-user.target
