# opt/config/includes.chroot/etc/systemd/system/tsig-rotate.timer
[Unit]
Description=Monthly TSIG key rotation

[Timer]
OnCalendar=monthly
Persistent=true
Unit=tsig-rotate.service

[Install]
WantedBy=timers.target
